# web-medium-php_sql-main/deploy/php/crontab
0 * * * * root /usr/bin/python3 /usr/local/bin/cleanup.py >> /var/log/cleanup.log 2>&1
